# infra/mdns/linux/celesteos-mdns.service
[Unit]
Description=CelesteOS mDNS Advertisement Service
After=network.target avahi-daemon.service
Wants=avahi-daemon.service

[Service]
Type=simple
User=nobody
WorkingDirectory=/opt/celesteos
ExecStart=/usr/bin/python3 /opt/celesteos/infra/mdns/advertise-mdns.py 80
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target